# ttmjserver110/webagent/ln.sh
#!/bin/sh

#chkconfig: - 85 15
#description: my game server

<<<<<<< HEAD
ln -s /root/webagent/public /root/llmpro/
=======
ln -s /root/webagent/public /root/webadmin/
>>>>>>> master
