# ttmjserver110/dataCenterService/statisticsTask/crontab/userConsumption.sh
#!/bin/bash

filename=${0}
extension="${filename##*.}"
currentFileName="${filename%.*}"

bash baseScript ${currentFileName}